# own_files/mt7628/files_mt7628shellgui/usr/shellgui/shellguilighttpd/www/apps/jimai/root.cron
*/10 * * * * /usr/shellgui/shellguilighttpd/www/apps/jimai/jimai.sbin check