# own_files/ar934x/files_zhonglian/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
/etc/gpioloop.sh &
#/etc/init.d/config4g restart
#sleep 25
/etc/modem_info.sh &
#/etc/plc/plc.pub &
/etc/plc/updateinfo.sh &
/etc/plc/mqtt_listen.sh &
/etc/init.d/n2n_v2 restart &
exit 0
