# own_files/files_mt7628shellgui/usr/shellgui/shellguilighttpd/www/apps/notice/root.cron
0 */2 * * * sleep 35; /usr/shellgui/progs/main.sbin sent_backup 1