# own_files/mt7628/files_wifiaudio_vlc/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

sleep 5;
shairport -a "wifiaudio_box" -P /tmp/airplay.pid --daemon;
sleep 5
#gmediarender -I 192.168.11.1 -f "linkall_box" -d;
/wifi_audio/bin/dlna --name=wifiaudio&
#echo 39 > /sys/class/gpio/export
#echo out > /sys/class/gpio/gpio39/direction
#echo 0 > /sys/class/gpio/gpio39/value
#sleep 1
echo 1 > /sys/class/gpio/gpio39/value


exit 0
