# own_files/qca953x/files_9531shellgui_4g_pppd/usr/shellgui/shellguilighttpd/www/apps/ac-set/root.cron
* * * * * /usr/shellgui/shellguilighttpd/www/apps/ac-set/S1200-ac.init keep_file
