# own_files/mt7628/files_mt7628shellgui/usr/shellgui/shellguilighttpd/www/apps/firewall-extra/root.cron
0 3 * * * /usr/shellgui/progs/main.sbin geoip_update
