# own_files/files_hi9344/etc/init.d/serialoverip
#!/bin/sh /etc/rc.common

START=99
STOP=99

start() {
 sleep 15
 /etc/start_serialoverip.sh &

}

stop() {
killall serialoverip
}
