# own_files/qca953x/files_9531miniv2/etc/crontabs/root
*/2 * * * *  /etc/setaponly
*/5 * * * * /etc/modem_info.sh
