# own_files/files_guide/etc/hotplug.d/usb/10-wirelss
if [ "$ACTION" = "add" ]; then
rm /etc/config/wireless
sleep 3
wifi detect > /etc/config/wireless
/etc/init.d/network restart
fi