# own_files/files_tw7628_wifidog/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

/etc/init.d/dnsmasq restart &
sleep 5 
/etc/init.d/network restart &
/etc/exportgpio.sh &
/etc/modem_info.sh &
/etc/watch_wifidog.sh &
iptables -N BLKLIST
iptables -I FORWARD -p tcp -j BLKLIST
exit 0
