# own_files/files_zhuguangyu/etc/crontabs/root
0 2 * * * reboot
*/5 * * * * /etc/rsync.sh
