# own_files/qca953x/files_9531shellgui_4g_pppd/usr/shellgui/shellguilighttpd/www/apps/qos-shellgui/root.cron
0 0,4,8,12,16,20 * * * date -u  +"%Y.%m.%d-%H:%M:%S" >/usr/data/time_backup
0,1,11,21,31,41,51 * * * * /usr/bin/s-k-t >/dev/null 2>&1

