# own_files/files_mt7628shellgui/usr/shellgui/shellguilighttpd/www/apps/lan-net-record/root.cron
* * * * * sleep 10;/usr/shellgui/shellguilighttpd/www/apps/lan-net-record/S1100-lan-net-record.init record
*/30 * * * * /usr/shellgui/progs/main.sbin socks_job