# own_files/qca953x/files_9531shellgui_gobinet_OK/etc/uci-defaults/10-init-4g-interface
#!/bin/sh
#
# Copyright (C) 2011 OpenWrt.org
#
uci set firewall.@zone[1].network="wan wan6 4g"
uci commit firewall

